# deploy/oracle/systemd/pv-premium-sync.service
[Unit]
Description=PV Premium Sync Backend
After=network.target

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/srv/pv-premium
EnvironmentFile=/etc/pv-premium-sync/backend.env
ExecStart=/usr/bin/node --experimental-sqlite backend/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/pv-premium-sync /var/backups/pv-premium-sync

[Install]
WantedBy=multi-user.target
